# system/ThirdParty/BazMessenger.conf
[program:BazMessenger]
command                 = bash -c "ulimit -n 10000; exec /usr/bin/php /var/www/html/sp/public/index.php messenger start"
process_name            = BazMessengerRatchet
numprocs                = 1
autostart               = true
autorestart             = true
user                    = root
stdout_logfile          = /var/www/html/sp/var/log/messenger-info.log
stdout_logfile_maxbytes = 1MB
stderr_logfile          = /var/www/html/sp/var/log/messenger-error.log
stderr_logfile_maxbytes = 1MB
